# renew-technitium-cert/etc/systemd/system/renew-technitium-cert.service
[Unit]
Description=Renew Technitium Let's Encrypt certificate
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/renew-technitium-cert.sh
