# pbs-file-backup/etc/systemd/system/pbs-file-backup@weekly.timer
[Unit]
Description=Monthly PBS filesystem backup (%i)

[Timer]
OnCalendar=Fri 03:00
RandomizedDelaySec=2h
Persistent=true

Unit=pbs-file-backup@weekly.service

[Install]
WantedBy=timers.target
